# systemd/pocketagent.service
[Unit]
Description=PocketAgent (voice reminder assistant)
After=network-online.target sound.target pocketagent-display.service
Requires=pocketagent-display.service
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/pocketagent
Environment=NODE_ENV=production
# Provide OPENAI_API_KEY via /etc/default/pocketagent or systemd drop-in
EnvironmentFile=-/etc/default/pocketagent
ExecStart=/usr/bin/node /opt/pocketagent/pocketagent/index.js

# Let reminders be a separate service; the chat agent stays responsive.
# (Port is configured via /etc/default/pocketagent)

Restart=on-failure
RestartSec=2
User=pi
Group=pi
SupplementaryGroups=audio gpio

[Install]
WantedBy=multi-user.target
