# systemd/pocketagent-display.service
[Unit]
Description=PocketAgent Display Sidecar (Whisplay LCD)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/pocketagent
Environment=NODE_ENV=production
EnvironmentFile=-/etc/default/pocketagent
ExecStart=/usr/bin/python3 /opt/pocketagent/whisplay/display/server.py
Restart=on-failure
RestartSec=2
User=pi
Group=pi

[Install]
WantedBy=multi-user.target
